# qr3_77.conf
[program:qr3_77]
command=python3 /home/curator3/qr3_77/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_77/test.err.log
stdout_logfile=/home/curator3/qr3_77/test.err.log
